# deploy/supervisor-worker.conf
[program:ai-inbox-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/todo/current/artisan queue:work --queue=default --tries=3 --sleep=1 --timeout=120
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/ai-inbox/worker.log
stopwaitsecs=120
